# service_files/auto-clean-logs.timer
[Unit]
Description=JAJA Auto Clean Logs Timer

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
