# service_files/auto-clean-logs.timer
[Unit]
Description=Daily Auto Clean Old Logs Timer

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
