# service_files/fedora-auto-setup.service
[Unit]
Description=JAJA Auto-Setup Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/02_full_auto_setup.sh --daemon
Restart=on-failure
RestartSec=30
EnvironmentFile=/etc/jaja.conf
TimeoutStartSec=300
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
